# config/systemd/nodeagent.service
[Unit]
Description=Zeitwork Node Agent
Documentation=https://github.com/zeitwork/zeitwork
After=network-online.target containerd.service docker.service
Wants=network-online.target
Requires=containerd.service

[Service]
Type=simple
User=root
Restart=always
RestartSec=10
EnvironmentFile=/etc/zeitwork/nodeagent.env
ExecStart=/usr/local/bin/nodeagent
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nodeagent

# Security
NoNewPrivileges=false
PrivateTmp=true

# Network - Allow binding to port 8080 for reverse proxy
# AmbientCapabilities=CAP_NET_BIND_SERVICE (not needed for ports > 1024)

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target

